# roles/harden_profile/files/etc/profile.d/tmout.sh
#!/bin/bash

TMOUT=600
readonly TMOUT
export TMOUT
